# roles/desktop_backup/templates/desktop-backup.timer.j2
[Unit]
Description=Daily user backup timer

[Timer]
OnCalendar=*-*-* {{ desktop_backup_time }}
Persistent=true

[Install]
WantedBy=timers.target
