# roles/docker_backup/templates/docker-backup.timer.j2
[Unit]
Description=Run Docker backup daily

[Timer]
OnCalendar={{ docker_backup_oncalendar }}
Persistent=true

[Install]
WantedBy=timers.target
