# roles/docker_backup/templates/docker-backup.service.j2
[Unit]
Description=Backup running Docker containers
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/docker-backup.sh
